# nomad-image/nomad-comfyui.service
[Unit]
Description=NOMAD local image generation (ComfyUI on the RTX 5090)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# %h = the user's home. Adjust if you cloned ComfyUI elsewhere (COMFYUI_DIR in install_comfyui.sh).
WorkingDirectory=%h/ComfyUI
# --listen 0.0.0.0 so the v2 engine container can reach it via host.docker.internal:8188.
# WSL2 is NAT'd (not on the LAN); ComfyUI has no auth, so keep this box's firewall closed.
# --disable-pinned-memory: big models (FLUX.2, 33GB DiT) offload to CPU RAM; pinned (page-locked,
#   non-swappable) memory can exceed WSL's RAM and hard-OOM ComfyUI. Pageable offload + swap is
#   slower but survives. --cache-none trims cached weights too. Needed on a 24GB card / 47GB WSL.
ExecStart=%h/ComfyUI/venv/bin/python %h/ComfyUI/main.py --listen 0.0.0.0 --port 8188 --disable-pinned-memory --cache-none
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
